# openclaw/rootfs/etc/s6-overlay/s6-rc.d/init-openclaw/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Initialize OpenClaw configuration
# ==============================================================================

bashio::log.info "Initializing OpenClaw addon..."

# Create config directories
mkdir -p /data/openclaw
mkdir -p "$(bashio::config 'workspace_path')"

# Get configuration values
ANTHROPIC_API_KEY=$(bashio::config 'anthropic_api_key')
HA_ACCESS_TOKEN=$(bashio::config 'ha_access_token')
LOG_LEVEL=$(bashio::config 'log_level')
WORKSPACE_PATH=$(bashio::config 'workspace_path')

# Determine Home Assistant access method
if bashio::config.has_value 'ha_access_token'; then
    bashio::log.info "Using provided Home Assistant access token"
    HA_TOKEN="${HA_ACCESS_TOKEN}"
    HA_URL="http://supervisor/core"
else
    bashio::log.info "Using Supervisor API token"
    HA_TOKEN="${SUPERVISOR_TOKEN}"
    HA_URL="http://supervisor/core"
fi

# Create OpenClaw config directory
OPENCLAW_HOME="/data/openclaw"
mkdir -p "${OPENCLAW_HOME}"

# Generate or reuse gateway token (persists across restarts)
TOKEN_FILE="${OPENCLAW_HOME}/.gateway-token"
if [[ -f "${TOKEN_FILE}" ]]; then
    GATEWAY_TOKEN=$(cat "${TOKEN_FILE}")
else
    GATEWAY_TOKEN=$(head -c 32 /dev/urandom | base64 | tr -dc 'a-zA-Z0-9' | head -c 32)
    echo "${GATEWAY_TOKEN}" > "${TOKEN_FILE}"
    chmod 600 "${TOKEN_FILE}"
fi

bashio::log.info "Gateway token configured"

# Log gateway token for user reference
bashio::log.info "Gateway token: ${GATEWAY_TOKEN}"

# Generate OpenClaw configuration in correct location
mkdir -p "${OPENCLAW_HOME}/.openclaw"
cat > "${OPENCLAW_HOME}/.openclaw/openclaw.json" << EOF
{
  "env": {
    "ANTHROPIC_API_KEY": "${ANTHROPIC_API_KEY}"
  },
  "agents": {
    "defaults": {
      "model": {
        "primary": "anthropic/claude-opus-4-5"
      }
    }
  },
  "gateway": {
    "port": 18789,
    "bind": "lan",
    "mode": "local",
    "trustedProxies": ["172.30.32.0/23", "127.0.0.1", "::1"],
    "auth": {
      "mode": "token",
      "token": "${GATEWAY_TOKEN}"
    },
    "controlUi": {
      "enabled": true,
      "allowInsecureAuth": true
    }
  },
  "workspace": "${WORKSPACE_PATH}",
  "logLevel": "${LOG_LEVEL}"
}
EOF

# Export environment variables for the main service
{
    echo "HA_TOKEN=${HA_TOKEN}"
    echo "HA_URL=${HA_URL}"
    echo "OPENCLAW_CONFIG_DIR=${OPENCLAW_HOME}"
    echo "OPENCLAW_WORKSPACE_DIR=${WORKSPACE_PATH}"
    echo "HOME=${OPENCLAW_HOME}"
    echo "OPENCLAW_GATEWAY_TOKEN=${GATEWAY_TOKEN}"
    echo "ANTHROPIC_API_KEY=${ANTHROPIC_API_KEY}"
} > /var/run/s6/container_environment/openclaw-env

# Install HA skill to OpenClaw
SKILL_DIR="${OPENCLAW_HOME}/skills/homeassistant"
mkdir -p "${SKILL_DIR}"
cp -r /opt/openclaw-ha-skill/* "${SKILL_DIR}/"

bashio::log.info "Home Assistant skill installed"

# Create auth-profiles.json for the main agent
AGENT_DIR="${OPENCLAW_HOME}/.openclaw/agents/main/agent"
mkdir -p "${AGENT_DIR}"

cat > "${AGENT_DIR}/auth-profiles.json" << EOF
{
  "version": 1,
  "profiles": {
    "anthropic:manual": {
      "type": "token",
      "provider": "anthropic",
      "token": "${ANTHROPIC_API_KEY}"
    }
  }
}
EOF

bashio::log.info "Anthropic API key configured"

bashio::log.info "OpenClaw initialization complete"
