# openclaw/rootfs/etc/s6-overlay/s6-rc.d/init-openclaw/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Initialize OpenClaw configuration
# ==============================================================================

bashio::log.info "Initializing OpenClaw addon..."

# Create config directories
mkdir -p /data/openclaw
mkdir -p "$(bashio::config 'workspace_path')"

# Get configuration values
ANTHROPIC_API_KEY=$(bashio::config 'anthropic_api_key')
HA_ACCESS_TOKEN=$(bashio::config 'ha_access_token')
LOG_LEVEL=$(bashio::config 'log_level')
WORKSPACE_PATH=$(bashio::config 'workspace_path')
SIGNAL_ENABLED=$(bashio::config 'signal_enabled')
SIGNAL_PHONE=$(bashio::config 'signal_phone')
VOICE_TRANSCRIPTION=$(bashio::config 'voice_transcription')
TTS_VOICE=$(bashio::config 'tts_voice')

# Determine Home Assistant access method
if bashio::config.has_value 'ha_access_token'; then
    bashio::log.info "Using provided Home Assistant access token"
    HA_TOKEN="${HA_ACCESS_TOKEN}"
    HA_URL="http://supervisor/core"
else
    bashio::log.info "Using Supervisor API token"
    HA_TOKEN="${SUPERVISOR_TOKEN}"
    HA_URL="http://supervisor/core"
fi

# Create OpenClaw config directory
OPENCLAW_HOME="/data/openclaw"
mkdir -p "${OPENCLAW_HOME}"

# Generate or reuse gateway token (persists across restarts)
TOKEN_FILE="${OPENCLAW_HOME}/.gateway-token"
if [[ -f "${TOKEN_FILE}" ]]; then
    GATEWAY_TOKEN=$(cat "${TOKEN_FILE}")
else
    GATEWAY_TOKEN=$(head -c 32 /dev/urandom | base64 | tr -dc 'a-zA-Z0-9' | head -c 32)
    echo "${GATEWAY_TOKEN}" > "${TOKEN_FILE}"
    chmod 600 "${TOKEN_FILE}"
fi

bashio::log.info "Gateway token configured"

# Log gateway token for user reference
bashio::log.info "Gateway token: ${GATEWAY_TOKEN}"

# Generate OpenClaw configuration in correct location
mkdir -p "${OPENCLAW_HOME}/.openclaw"

# Build Signal channel config if enabled
SIGNAL_CHANNEL_CONFIG=""
if bashio::var.true "${SIGNAL_ENABLED}"; then
    bashio::log.info "Signal channel enabled for phone: ${SIGNAL_PHONE}"
    mkdir -p "${OPENCLAW_HOME}/.signal-cli"
    SIGNAL_CHANNEL_CONFIG=$(cat << SIGNAL_EOF
  "channels": {
    "signal": {
      "enabled": true,
      "account": "${SIGNAL_PHONE}",
      "cliPath": "signal-cli",
      "dmPolicy": "open",
      "allowFrom": ["*"]
    }
  },
SIGNAL_EOF
)
fi

# Build audio transcription config if enabled
AUDIO_CONFIG=""
if [[ "${VOICE_TRANSCRIPTION}" == "sherpa-onnx" ]]; then
    bashio::log.info "Voice transcription enabled: sherpa-onnx"
    # Use wrapper script that handles AAC->WAV conversion
    if [[ -x "/usr/local/bin/transcribe-audio" ]]; then
        bashio::log.info "Using transcribe-audio wrapper (handles AAC->WAV conversion)"
        AUDIO_CONFIG=$(cat << AUDIO_EOF
  "tools": {
    "media": {
      "audio": {
        "enabled": true,
        "models": [
          {
            "type": "cli",
            "command": "/usr/local/bin/transcribe-audio",
            "args": ["{{MediaPath}}"]
          }
        ]
      }
    }
  },
AUDIO_EOF
)
    else
        bashio::log.warning "transcribe-audio wrapper not found"
    fi
fi

# Build TTS config for German voice (Edge TTS)
TTS_CONFIG=""
if bashio::config.has_value 'tts_voice' || [[ -n "${TTS_VOICE}" ]]; then
    TTS_VOICE_FINAL="${TTS_VOICE:-de-DE-ConradNeural}"
    bashio::log.info "TTS configured: Edge with voice ${TTS_VOICE_FINAL}"
    TTS_CONFIG=$(cat << TTS_EOF
  "messages": {
    "tts": {
      "auto": "always",
      "provider": "edge",
      "edge": {
        "voice": "${TTS_VOICE_FINAL}"
      }
    }
  },
TTS_EOF
)
fi

cat > "${OPENCLAW_HOME}/.openclaw/openclaw.json" << EOF
{
  "env": {
    "ANTHROPIC_API_KEY": "${ANTHROPIC_API_KEY}"
  },
  "agents": {
    "defaults": {
      "model": {
        "primary": "anthropic/claude-opus-4-5"
      },
      "contextPruning": {
        "mode": "cache-ttl",
        "ttl": "30m",
        "keepLastAssistants": 5
      }
    }
  },
${AUDIO_CONFIG}
${TTS_CONFIG}
${SIGNAL_CHANNEL_CONFIG}
  "gateway": {
    "port": 18789,
    "bind": "lan",
    "mode": "local",
    "trustedProxies": ["172.30.32.0/23", "127.0.0.1", "::1"],
    "auth": {
      "mode": "token",
      "token": "${GATEWAY_TOKEN}"
    },
    "controlUi": {
      "enabled": true,
      "allowInsecureAuth": true
    }
  }
}
EOF

# Export environment variables for the main service
{
    echo "HA_TOKEN=${HA_TOKEN}"
    echo "HA_URL=${HA_URL}"
    echo "OPENCLAW_CONFIG_DIR=${OPENCLAW_HOME}"
    echo "OPENCLAW_WORKSPACE_DIR=${WORKSPACE_PATH}"
    echo "HOME=${OPENCLAW_HOME}"
    echo "OPENCLAW_GATEWAY_TOKEN=${GATEWAY_TOKEN}"
    echo "ANTHROPIC_API_KEY=${ANTHROPIC_API_KEY}"
    # Signal settings via env vars
    echo "SIGNAL_CLI_CONFIG_DIR=${OPENCLAW_HOME}/.signal-cli"
    echo "SIGNAL_PHONE=${SIGNAL_PHONE}"
    echo "SIGNAL_TRANSCRIPTION=${VOICE_TRANSCRIPTION}"
} > /var/run/s6/container_environment/openclaw-env

# Install HA skill to OpenClaw
SKILL_DIR="${OPENCLAW_HOME}/skills/homeassistant"
mkdir -p "${SKILL_DIR}"
cp -r /opt/openclaw-ha-skill/* "${SKILL_DIR}/"

bashio::log.info "Home Assistant skill installed"

# Create auth-profiles.json for the main agent
AGENT_DIR="${OPENCLAW_HOME}/.openclaw/agents/main/agent"
mkdir -p "${AGENT_DIR}"

cat > "${AGENT_DIR}/auth-profiles.json" << EOF
{
  "version": 1,
  "profiles": {
    "anthropic:manual": {
      "type": "token",
      "provider": "anthropic",
      "token": "${ANTHROPIC_API_KEY}"
    }
  }
}
EOF

bashio::log.info "Anthropic API key configured"

# Run openclaw doctor --fix if Signal is enabled to apply channel changes
if bashio::var.true "${SIGNAL_ENABLED}"; then
    bashio::log.info "Running openclaw doctor --fix to enable Signal channel..."
    cd "${WORKSPACE_PATH}" || exit 1
    HOME="${OPENCLAW_HOME}" openclaw doctor --fix || bashio::log.warning "openclaw doctor --fix returned non-zero"
fi

bashio::log.info "OpenClaw initialization complete"
