# openclaw/rootfs/etc/s6-overlay/s6-rc.d/init-openclaw/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Initialize OpenClaw configuration
# ==============================================================================

bashio::log.info "Initializing OpenClaw addon..."

# Create config directories
mkdir -p /data/openclaw
mkdir -p "$(bashio::config 'workspace_path')"

# Apply DNS/hosts overrides
MARKER="# OpenClaw DNS Override"
# Remove old OpenClaw entries first
sed -i "/${MARKER}/d" /etc/hosts 2>/dev/null || true

DNS_COUNT=$(bashio::config 'dns_overrides | length')
if [[ "${DNS_COUNT}" -gt 0 ]]; then
    bashio::log.info "Applying ${DNS_COUNT} DNS override(s)..."
    for ((i=0; i<DNS_COUNT; i++)); do
        ENTRY=$(bashio::config "dns_overrides[${i}]")
        if [[ -n "${ENTRY}" ]]; then
            echo "${ENTRY} ${MARKER}" >> /etc/hosts
            bashio::log.info "  Added: ${ENTRY}"
        fi
    done
fi

# Get configuration values
ANTHROPIC_API_KEY=$(bashio::config 'anthropic_api_key')
HA_ACCESS_TOKEN=$(bashio::config 'ha_access_token')
LOG_LEVEL=$(bashio::config 'log_level')
WORKSPACE_PATH=$(bashio::config 'workspace_path')
SIGNAL_ENABLED=$(bashio::config 'signal_enabled')
SIGNAL_PHONE=$(bashio::config 'signal_phone')
VOICE_TRANSCRIPTION=$(bashio::config 'voice_transcription')
TTS_VOICE=$(bashio::config 'tts_voice')

# Determine Home Assistant access method
if bashio::config.has_value 'ha_access_token'; then
    bashio::log.info "Using provided Home Assistant access token"
    HA_TOKEN="${HA_ACCESS_TOKEN}"
    HA_URL="http://supervisor/core"
else
    bashio::log.info "Using Supervisor API token"
    HA_TOKEN="${SUPERVISOR_TOKEN}"
    HA_URL="http://supervisor/core"
fi

# Create OpenClaw config directory
OPENCLAW_HOME="/data/openclaw"
mkdir -p "${OPENCLAW_HOME}"

# Generate or reuse gateway token (persists across restarts)
TOKEN_FILE="${OPENCLAW_HOME}/.gateway-token"
if [[ -f "${TOKEN_FILE}" ]]; then
    GATEWAY_TOKEN=$(cat "${TOKEN_FILE}")
else
    GATEWAY_TOKEN=$(head -c 32 /dev/urandom | base64 | tr -dc 'a-zA-Z0-9' | head -c 32)
    echo "${GATEWAY_TOKEN}" > "${TOKEN_FILE}"
    chmod 600 "${TOKEN_FILE}"
fi

bashio::log.info "Gateway token configured"

# Log gateway token for user reference
bashio::log.info "Gateway token: ${GATEWAY_TOKEN}"

# Generate OpenClaw configuration in correct location
mkdir -p "${OPENCLAW_HOME}/.openclaw"

# Build Signal channel config if enabled
SIGNAL_CHANNEL_CONFIG=""
if bashio::var.true "${SIGNAL_ENABLED}"; then
    bashio::log.info "Signal channel enabled for phone: ${SIGNAL_PHONE}"
    mkdir -p "${OPENCLAW_HOME}/.signal-cli"
    SIGNAL_CHANNEL_CONFIG=$(cat << SIGNAL_EOF
  "channels": {
    "signal": {
      "enabled": true,
      "account": "${SIGNAL_PHONE}",
      "cliPath": "signal-cli",
      "dmPolicy": "open",
      "allowFrom": ["*"]
    }
  },
SIGNAL_EOF
)
fi

# Build audio transcription config if enabled
AUDIO_CONFIG=""
if [[ "${VOICE_TRANSCRIPTION}" == "sherpa-onnx" ]]; then
    bashio::log.info "Voice transcription enabled: sherpa-onnx"
    # Use wrapper script that handles AAC->WAV conversion
    if [[ -x "/usr/local/bin/transcribe-audio" ]]; then
        bashio::log.info "Using transcribe-audio wrapper (handles AAC->WAV conversion)"
        AUDIO_CONFIG=$(cat << AUDIO_EOF
  "tools": {
    "media": {
      "audio": {
        "enabled": true,
        "models": [
          {
            "type": "cli",
            "command": "/usr/local/bin/transcribe-audio",
            "args": ["{{MediaPath}}"]
          }
        ]
      }
    }
  },
AUDIO_EOF
)
    else
        bashio::log.warning "transcribe-audio wrapper not found"
    fi
fi

# Build TTS config for German voice (Edge TTS)
TTS_CONFIG=""
if bashio::config.has_value 'tts_voice' || [[ -n "${TTS_VOICE}" ]]; then
    TTS_VOICE_FINAL="${TTS_VOICE:-de-DE-ConradNeural}"
    bashio::log.info "TTS configured: Edge with voice ${TTS_VOICE_FINAL}"
    TTS_CONFIG=$(cat << TTS_EOF
  "messages": {
    "tts": {
      "auto": "always",
      "provider": "edge",
      "edge": {
        "voice": "${TTS_VOICE_FINAL}"
      }
    }
  },
TTS_EOF
)
fi

cat > "${OPENCLAW_HOME}/.openclaw/openclaw.json" << EOF
{
  "env": {
    "ANTHROPIC_API_KEY": "${ANTHROPIC_API_KEY}"
  },
  "agents": {
    "defaults": {
      "workspace": "${WORKSPACE_PATH}",
      "model": {
        "primary": "anthropic/claude-opus-4-5"
      },
      "contextPruning": {
        "mode": "cache-ttl",
        "ttl": "15m",
        "keepLastAssistants": 2
      }
    }
  },
${AUDIO_CONFIG}
${TTS_CONFIG}
${SIGNAL_CHANNEL_CONFIG}
  "gateway": {
    "port": 18789,
    "bind": "lan",
    "mode": "local",
    "trustedProxies": ["172.30.32.0/23", "127.0.0.1", "::1"],
    "auth": {
      "mode": "token",
      "token": "${GATEWAY_TOKEN}"
    },
    "controlUi": {
      "enabled": true,
      "allowInsecureAuth": true
    }
  }
}
EOF

# Export environment variables for the main service
{
    echo "HA_TOKEN=${HA_TOKEN}"
    echo "HA_URL=${HA_URL}"
    echo "OPENCLAW_CONFIG_DIR=${OPENCLAW_HOME}"
    echo "OPENCLAW_WORKSPACE_DIR=${WORKSPACE_PATH}"
    echo "HOME=${OPENCLAW_HOME}"
    echo "OPENCLAW_GATEWAY_TOKEN=${GATEWAY_TOKEN}"
    echo "ANTHROPIC_API_KEY=${ANTHROPIC_API_KEY}"
    # Signal settings via env vars
    echo "SIGNAL_CLI_CONFIG_DIR=${OPENCLAW_HOME}/.signal-cli"
    echo "SIGNAL_PHONE=${SIGNAL_PHONE}"
    echo "SIGNAL_TRANSCRIPTION=${VOICE_TRANSCRIPTION}"
} > /var/run/s6/container_environment/openclaw-env

# Install HA skill to OpenClaw
SKILL_DIR="${OPENCLAW_HOME}/skills/homeassistant"
mkdir -p "${SKILL_DIR}"
cp -r /opt/openclaw-ha-skill/* "${SKILL_DIR}/"

bashio::log.info "Home Assistant skill installed"

# Create workspace bootstrap files (only if missing - preserves user customizations)
create_if_missing() {
    local file="$1"
    local content="$2"
    if [[ ! -f "${WORKSPACE_PATH}/${file}" ]]; then
        echo "${content}" > "${WORKSPACE_PATH}/${file}"
        bashio::log.info "Created ${file}"
    else
        bashio::log.debug "${file} already exists, skipping"
    fi
}

# AGENTS.md - Operating instructions and memory
create_if_missing "AGENTS.md" "# Home Assistant Voice Assistant

You are a Home Assistant voice assistant running via OpenClaw. Your primary purpose is controlling smart home devices.

## Available Tools

You have access to the \`homeassistant\` skill with these tools:
- \`ha_list_devices\` - List all available devices (use sparingly, cache results mentally)
- \`ha_get_states\` - Get all entity states
- \`ha_get_state\` - Get specific entity state
- \`ha_turn_on\` / \`ha_turn_off\` / \`ha_toggle\` - Control devices
- \`ha_set_brightness\` - Set light brightness (0-100)
- \`ha_set_temperature\` - Set thermostat temperature
- \`ha_trigger_automation\` - Trigger automations
- \`ha_run_script\` - Run HA scripts

## Operating Guidelines

1. **Be concise** - Voice responses should be short and natural
2. **Confirm actions** - \"Turning on the living room light\"
3. **Handle errors gracefully** - \"I couldn't find that device\"
4. **Remember context** - If user says \"turn it off\", know what \"it\" refers to
5. **Don't list all devices** - Only call ha_list_devices when explicitly needed"

# SOUL.md - Persona and behavioral boundaries
create_if_missing "SOUL.md" "# Personality

You are a helpful, efficient home assistant. You speak naturally and concisely, like a smart speaker assistant.

## Communication Style
- Brief, conversational responses
- German language (match user's language)
- No unnecessary explanations
- Confirm actions taken
- Ask for clarification when ambiguous

## Boundaries
- Only control devices the user mentions
- Don't make assumptions about device states
- Don't execute potentially dangerous actions without confirmation"

# TOOLS.md - Tool usage notes and patterns
create_if_missing "TOOLS.md" "# Tool Usage Guidelines

## Home Assistant Skill

**IMPORTANT:** Avoid calling \`ha_list_devices\` or \`ha_get_states\` repeatedly. These return hundreds of entities and bloat context.

### Efficient Patterns
- For specific devices: Use \`ha_get_state\` with known entity_id
- For control: Go directly to \`ha_turn_on\`/\`ha_turn_off\` with entity_id
- For discovery: Only use \`ha_list_devices\` when user explicitly asks \"what devices do I have?\"

### Entity ID Patterns
Common domains: \`light.\`, \`switch.\`, \`climate.\`, \`cover.\`, \`script.\`, \`automation.\`"

# IDENTITY.md - Agent name and character
create_if_missing "IDENTITY.md" "name: Assistant
emoji: ðŸ 
character: Home automation assistant"

# USER.md - User profile and preferences
create_if_missing "USER.md" "# User Profile

Language: German
Location: Home
Interaction style: Voice commands via Signal/HA Assist"

bashio::log.info "Workspace bootstrap files configured"

# Create auth-profiles.json for the main agent
AGENT_DIR="${OPENCLAW_HOME}/.openclaw/agents/main/agent"
mkdir -p "${AGENT_DIR}"

cat > "${AGENT_DIR}/auth-profiles.json" << EOF
{
  "version": 1,
  "profiles": {
    "anthropic:manual": {
      "type": "token",
      "provider": "anthropic",
      "token": "${ANTHROPIC_API_KEY}"
    }
  }
}
EOF

bashio::log.info "Anthropic API key configured"

# Run openclaw doctor --fix if Signal is enabled to apply channel changes
if bashio::var.true "${SIGNAL_ENABLED}"; then
    bashio::log.info "Running openclaw doctor --fix to enable Signal channel..."
    cd "${WORKSPACE_PATH}" || exit 1
    HOME="${OPENCLAW_HOME}" openclaw doctor --fix || bashio::log.warning "openclaw doctor --fix returned non-zero"
fi

# Run openclaw setup to initialize workspace properly
bashio::log.info "Running openclaw setup to initialize workspace..."
cd "${WORKSPACE_PATH}" || exit 1
HOME="${OPENCLAW_HOME}" openclaw setup --workspace "${WORKSPACE_PATH}" 2>/dev/null || bashio::log.debug "openclaw setup completed (or already initialized)"

bashio::log.info "OpenClaw initialization complete"
