# openclaw/rootfs/etc/s6-overlay/s6-rc.d/signal-cli/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
bashio::log.info "signal-cli daemon stopped"
