# openclaw/rootfs/etc/s6-overlay/s6-rc.d/signal-cli/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start signal-cli daemon for OpenClaw Signal channel
# ==============================================================================

# Source environment variables
if [[ -f /var/run/s6/container_environment/openclaw-env ]]; then
    # shellcheck disable=SC1091
    source /var/run/s6/container_environment/openclaw-env
fi

# Check if Signal is enabled
SIGNAL_ENABLED=$(bashio::config 'signal_enabled')
if ! bashio::var.true "${SIGNAL_ENABLED}"; then
    bashio::log.info "Signal channel disabled, not starting signal-cli daemon"
    exec sleep infinity
fi

# Set signal-cli config directory
export SIGNAL_CLI_CONFIG_DIR="${SIGNAL_CLI_CONFIG_DIR:-/data/openclaw/.signal-cli}"

# Wait for an account to be linked
ACCOUNTS_FILE="${SIGNAL_CLI_CONFIG_DIR}/data/accounts.json"

while true; do
    if [[ -f "${ACCOUNTS_FILE}" ]]; then
        # Extract first registered phone number from accounts.json
        # Format: {"accounts":[{"number":"+1234567890",...}]}
        REGISTERED_PHONE=$(grep -oP '"\+[0-9]+' "${ACCOUNTS_FILE}" | head -1 | tr -d '"')

        if [[ -n "${REGISTERED_PHONE}" ]]; then
            bashio::log.info "Found registered account: ${REGISTERED_PHONE}"
            bashio::log.info "Starting signal-cli daemon on HTTP port 8080..."
            exec signal-cli --config "${SIGNAL_CLI_CONFIG_DIR}" -a "${REGISTERED_PHONE}" daemon --http=127.0.0.1:8080
        fi
    fi

    bashio::log.info "Waiting for Signal account to be linked via web UI (port 18790)..."
    sleep 10
done
