# openclaw/rootfs/etc/s6-overlay/s6-rc.d/openclaw/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start OpenClaw gateway
# ==============================================================================

bashio::log.info "Starting OpenClaw gateway..."

# Source environment variables
if [[ -f /var/run/s6/container_environment/openclaw-env ]]; then
    # shellcheck disable=SC1091
    source /var/run/s6/container_environment/openclaw-env
fi

# Set HOME for OpenClaw
export HOME="${OPENCLAW_CONFIG_DIR:-/data/openclaw}"
export NODE_ENV="production"

# Change to workspace directory
cd "${OPENCLAW_WORKSPACE_DIR:-/share/openclaw}" || exit 1

# Read gateway token
GATEWAY_TOKEN=""
if [[ -f "${HOME}/.gateway-token" ]]; then
    GATEWAY_TOKEN=$(cat "${HOME}/.gateway-token")
fi

# Start OpenClaw gateway
if [[ -n "${GATEWAY_TOKEN}" ]]; then
    exec openclaw gateway --port 18789 --bind lan --allow-unconfigured --token "${GATEWAY_TOKEN}"
else
    bashio::log.error "Gateway token not found!"
    exit 1
fi
